# etc/profile.d/java_window_parenting.sh
#!/bin/sh
# This script will solve issues with Java windows on tiling WM such as bspwm.
export _JAVA_AWT_WM_NONREPARENTING=1
