# etc/profile.d/steam.sh
#!/bin/sh
# Steam will hang every time you try to close the main window.
# Just export this so it does it apropriately.
export STEAM_FRAME_FORCE_CLOSE=1
